# tomcat/tomcat.conf
[program:tomcat7]
command=/bin/bash -c "/tomcat/run.sh"
redirect_stderr=true
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
